# environment/src/main/bin/spark2/spark-slave/slave.conf
[supervisord]
nodaemon=true

[program:spark-slave]
command=sbin/start-slave.sh spark://master:7077
