# environment/src/main/bin/spark2/spark-master/master.conf
[supervisord]
nodaemon=true

[program:spark-master]
command=sbin/start-master.sh
